# rootfs/etc/s6-overlay/s6-rc.d/aria2/run
#!/command/execlineb -P

aria2c --conf-path=/aria2/config/aria2.conf