# rootfs/etc/s6-overlay/s6-rc.d/aria2-init/up
/etc/s6-overlay/scripts/setup
/etc/s6-overlay/scripts/init