# rootfs/etc/s6-overlay/s6-rc.d/aria2-cron/finish
#!/usr/bin/execlineb -S0

foreground { s6-test "$1" -eq 125 } exit 125